# get_pw.cron
*/30 * * * * env - `cat /etc/environment` /base/get_pw >> /dev/null 2>&1

